# sonicadrone/sonicadrone.conf
[program:sonicadrone]
directory=/home/rif
command=/home/rif/sonicadrone
autostart=true
autorestart=true
user=rif
redirect_stderr=true
